# sample/bbsd.service.in
[Unit]
Description=DreamBBS Telnet server
After=syslog.target network.target remote-fs.target nss-lookup.target startbbs.service
Requires=startbbs.service
Wants=xchatd.service

[Service]
Type=forking
PIDFile=@BBSHOME@/run/bbs.pid
ExecStart=@BBSHOME@/bin/bbsd 23
ExecStop=/bin/kill -s TERM $MAINPID
Restart=on-failure
RestartSec=42s
KillMode=process
PrivateTmp=true

[Install]
WantedBy=multi-user.target
