# sample/crontab
#!/bin/sh
# minute hour monthday month weekday command

# 系統報時
#5 8-17 * * 1-6 exec bin/classtable_alert > /dev/null 2>&1
#25 18-20 * * 1-6 exec bin/classtable_alert > /dev/null 2>&1

# 每天備份一次 source
#0 5 * * * /home/bbs/backsrc.sh > /dev/null 2>&1

#0 8 * * * /home/bbs/checkbrdDIR.sh

# 每小時轉信出去 2 次
#8,18,28,38,48,58 * * * * exec innd/bbslink /home/bbs > /dev/null 2>&1
#22,46 * * * * innd/bbslink > /dev/null 2>&1

# 每小時作一次人次統計及開票
5 */3 * * * exec bin/account > /dev/null 2>&1

# 看板資訊統計
#40 * * * * bin/brdstat > /dev/null 2>&1

# 對使用者信件作 expire
#20 5 4 * * exec bin/mailexpire > /dev/null 2>&1

# 每天做生日年齡統計
#0 10 * * * exec bin/countbirth > /dev/null 2>&1
#20 10 * * * exec bin/countage > /dev/null 2>&1
#40 10 * * * exec bin/countstar > /dev/null 2>&1

# 更新各版面的檔信列表
#50 */2 * * * bin/makefw > /dev/null 2>&1

# 每星期三增加點歌次數
40 4 * * 3 exec bin/addsong 4 1 > /dev/null 2>&1 #usage: bin/addsong [權限編號(perm)] [增加次數(number)]

# 每天清理個人信箱, 風之塔改成 quota 制
#10 5 * * * exec bin/bquota > /dev/null 2>&1

# 每小時更新一次動態看版
5 * * * * exec bin/acpro > /dev/null 2>&1
43 * * * * exec bin/camera > /dev/null 2>&1

# 每小時作一次話題統計, 特別注意是在 account執行之前
50 * * * * exec bin/poststat > /dev/null 2>&1

# 每天根據 expire.conf對看版文章作expire
20 2 * * * exec bin/expire > /dev/null 2>&1

# 每天對精華區作 garbage collection, garbage的形成詳見精華區編輯之document
30 5 * * * exec bin/gem-check > /dev/null 2>&1

#每天對精華區作索引檔, 方便使用者搜尋
20 6 * * 1,4 exec bin/gem-index > /dev/null 2>&1

# 每小時作歷史紀錄
58 * * * * exec bin/counter > /dev/null 2>&1

# 每星期一清除看板紀錄
20 7 * * 1 /home/bbs/sh/clean.sh > /dev/null 2>&1

# 每天掃描一次使用者的 ACCT
30 7 * * * /home/bbs/sh/checkACCT.sh > /dev/null 2>&1

# 每天掃描一次使用者信箱的 .DIR
35 7 * * * /home/bbs/sh/checkusrDIR.sh > /dev/null 2>&1

# 每天掃描一次看板的 .DIR
40 7 * * * /home/bbs/sh/checkbrdDIR.sh > /dev/null 2>&1

# 每天更新註冊資料
0 5 * * * exec bin/checkemail > /dev/null 2>&1

# 排行榜
0 4,10,15,20 * * * /home/bbs/sh/top.sh

# comment out for summer vacation
# 每天清理久未上站使用者, 清出的user會被搬至 usr/@目錄下, 請自行再刪除
#30 5 * * * exec bin/reaper > /dev/null 2>&1

# 每天清理個人信箱, 風之塔改成 quota 制
#10 5 * * * exec bin/bquota > /dev/null 2>&1

# 每十分鐘清除一次佔用系統資源過大之 bbsd
*/6 * * * * /home/bbs/sh/killtop.sh > /dev/null 2>&1

# 每天排定的備份
#0 4 * * * exec /home/bbs/bin/backup > /dev/null 2>&1

# 每天備份gem/@
#0 4 * * * /home/bbs/sh/backupgem@.sh > /dev/null 2>&1

# 每週四早上 6:15 定時重開機
#15 6 * * 4 /sbin/shutdown -r now > /dev/null 2>&1

