# sample/startbbs.service.in
[Unit]
Description=Start DreamBBS Telnet server
After=syslog.target network.target remote-fs.target nss-lookup.target

[Service]
User=@BBSUSR@
Group=@BBSGROUP@
Type=oneshot
## if set User and Group as bbs, then we don't need sudo commands
ExecStart=@BBSHOME@/bin/camera
ExecStart=@BBSHOME@/bin/account
ExecStart=@BBSHOME@/bin/acpro
ExecStart=@BBSHOME@/bin/makefw
PrivateTmp=true

[Install]
WantedBy=multi-user.target
